# rc.local.example
#!/bin/sh
sudo mkdir /tmp/qrencode
sudo mount -t tmpfs node /tmp/qrencode
